# deployment/systemctl/dev.socket
[Unit]
Description=dev gunicorn socket

[Socket]
ListenStream=/run/dev-gunicorn.sock

[Install]
WantedBy=sockets.target
