# deployment/systemctl/production.socket
[Unit]
Description=production gunicorn socket

[Socket]
ListenStream=/run/production-gunicorn.sock

[Install]
WantedBy=sockets.target