# deployment/systemctl/staging.socket
[Unit]
Description=staging gunicorn socket

[Socket]
ListenStream=/run/staging-gunicorn.sock

[Install]
WantedBy=sockets.target
